# autobuild-ceph.conf
description "autobuild ceph"
author "Tommi Virtanenh <tv@hq.newdream.net>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

nice 10
chdir /srv/autobuild-ceph

exec ./run
